# deploy/todayfood-backup.service
[Unit]
Description=Backup TodayFood logs and admin database (snowflow-site)

[Service]
Type=oneshot
WorkingDirectory=/home/ospacer/snowflow-site/todayfood
Environment=PROJECT_ROOT=/home/ospacer/snowflow-site/todayfood
ExecStart=/usr/bin/env bash /home/ospacer/snowflow-site/todayfood/deploy/backup_logs.sh
