# deploy/todayfood-admin.service
[Unit]
Description=今日宜吃 / today food Admin FastAPI App (snowflow-site, optional standalone)
After=network.target todayfood.service

[Service]
Type=simple
WorkingDirectory=/home/ospacer/snowflow-site/todayfood
Environment=HOST=127.0.0.1
Environment=PORT=9082
Environment=MONITORED_HEALTH_URL=http://127.0.0.1:9081/health
EnvironmentFile=/home/ospacer/snowflow-site/todayfood/backend/.env
ExecStart=/home/ospacer/.py_ai/bin/python -m uvicorn admin.main:app --host 127.0.0.1 --port 9082
Restart=always
RestartSec=5
StandardOutput=append:/home/ospacer/snowflow-site/todayfood/logs/admin.log
StandardError=append:/home/ospacer/snowflow-site/todayfood/logs/admin.log

[Install]
WantedBy=default.target
