# deploy/todayfood-backup.timer
[Unit]
Description=Run 今日宜吃 / today food backup daily

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
Unit=todayfood-backup.service

[Install]
WantedBy=timers.target
